# app-controller/systemd/ai-controller.service
[Unit]
Description=AI Controller Python Service
After=network.target redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/AIClient-2-API/app-controller
Environment=PATH=/root/AIClient-2-API/app-controller/.venv/bin
ExecStart=/root/AIClient-2-API/app-controller/.venv/bin/python main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiclient-python

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/root/AIClient-2-API/app-controller/logs

[Install]
WantedBy=multi-user.target
